# etc/local.d/hid_apple.start
#!/bin/bash

echo 1 > /sys/module/hid_apple/parameters/fnmode

echo 0 > /sys/module/hid_apple/parameters/iso_layout

echo 1 > /sys/module/hid_apple/parameters/swap_opt_cmd
